# ccnet.service
[Unit]
Description=Seafile Ccnet
After=network.target

[Service]
LimitNOFILE=10000
EnvironmentFile=/etc/sysconfig/seafile
Type=forking
PIDFile=/run/seafile/ccnet-server.pid
GuessMainPID=no
User=seafile
Group=seafile
RuntimeDirectoryPreserve=yes
RuntimeDirectory=seafile
ExecStart=/usr/sbin/ddns-notify
TimeoutStartSec=6
TimeoutStopSec=30


[Install]
WantedBy=multi-user.target