# seafile.service
[Unit]
Description=Seafile
Requires=ccnet.target
After=ccnet.target

[Service]
LimitNOFILE=10000
EnvironmentFile=/etc/sysconfig/seafile
Type=forking
PIDFile=/run/seafile-server.pid
GuessMainPID=no
User=seafile
Group=seafile
ExecStart=/usr/bin/seafutil-luanch --action seafile --config ${CONFIGFILE} -pidfile /run/seafile-server.pid -t 5
TimeoutStartSec=6
TimeoutStopSec=30


[Install]
WantedBy=multi-user.target