# seafile.service
[Unit]
Description=Seafile Server
Requires=ccnet.service
After=ccnet.service postgresql-15.service

[Service]
LimitNOFILE=10000
Type=simple
User=seafile
Group=seafile
ExecStart=/usr/bin/seafile-launch seafile --pidfile /run/seafile/seafile-server.pid
TimeoutStartSec=5
TimeoutStopSec=30


[Install]
WantedBy=multi-user.target
